# .deploy/config/ws.spinningcube.co.nz.conf
[program:ws.spinningcube]
command=/srv/app/artisan websockets:serve
autostart=true
autorestart=true
